# etc/service/70resque_web/run
#!/bin/bash
[ -r /etc/default/diaspora ] && . /etc/default/diaspora
cd $DIASPORA_HOME && HOME=/root exec /usr/local/bin/bundle exec resque-web -F -o localhost -p $RESQUE_WEB_PORT
